# docker-config/supervisor-app.conf
[program:app-uwsgi]
environment = DJANGO_SETTINGS_MODULE="podd.settings.docker", LANG=en_US.utf8, LC_ALL=en_US.utf8
command = /usr/local/bin/uwsgi --ini /usr/config/uwsgi.ini
